# etc/update-motd.d/10-hydra-status
#!/bin/bash

GREEN='\033[0;32m'
YELLOW='\033[1;33m'
RED='\033[0;31m'
BOLD='\033[1m'
NC='\033[0m'

UPTIME=$(uptime -p | cut -d " " -f 2-)
KERNEL=$(uname -r)
IP_ADDR=$(hostname -I | awk '{print $1}')
LOAD=$(cut -d ' ' -f1 /proc/loadavg)
MEM_TOTAL=$(awk '/MemTotal/ {printf "%.0f", $2/1024}' /proc/meminfo)
MEM_FREE=$(awk '/MemAvailable/ {printf "%.0f", $2/1024}' /proc/meminfo)
MEM_USED=$((MEM_TOTAL - MEM_FREE))

echo ""
echo -e "${BOLD}Project Hydra  —  NRF24-Gateway-01${NC}"
echo -e "RF-Hub / Data Relay Node"
echo ""
echo -e "${YELLOW}System${NC}"
echo -e "  Uptime:    ${UPTIME}"
echo -e "  Kernel:    ${KERNEL}"
echo -e "  IP:        ${IP_ADDR}"
echo -e "  Load:      ${LOAD}"
echo -e "  Memory:    ${MEM_USED} MB used / ${MEM_TOTAL} MB total"
echo ""
echo -e "${YELLOW}Hardware${NC}"
echo -e "  CPU ............ [ ${GREEN}OK${NC} ]"
echo -e "  RAM ............ [ ${GREEN}OK${NC} ]"
echo -e "  Storage ........ [ ${GREEN}OK${NC} ]"
echo -e "  Network ........ [ ${GREEN}OK${NC} ]"
echo -e "  RF-Module (NRF24L01+) ... [ ${GREEN}SYNC ACTIVE${NC} ]"
echo ""
echo -e "${YELLOW}Services${NC}"

# Python3 vorhanden?
if command -v python3 >/dev/null 2>&1; then
    PY_VER=$(python3 --version 2>&1 | awk '{print $2}')
    echo -e "  python3 ........ [ ${GREEN}${PY_VER}${NC} ]"
else
    echo -e "  python3 ........ [ ${RED}NOT FOUND${NC} ]"
fi

# RF-Skript vorhanden?
if [ -f /home/finn/hydra/hub/send_test.py ]; then
    echo -e "  hub/send_test .. [ ${GREEN}PRESENT${NC} ]"
else
    echo -e "  hub/send_test .. [ ${RED}MISSING${NC} ]"
fi

# SPI-Interface aktiv?
if ls /dev/spidev* >/dev/null 2>&1; then
    echo -e "  SPI interface .. [ ${GREEN}ACTIVE${NC} ]"
else
    echo -e "  SPI interface .. [ ${YELLOW}DISABLED${NC} ]"
fi

echo ""
echo -e "${YELLOW}Network${NC}"
/usr/local/bin/check_backend.sh
echo ""
